# root/etc/s6-overlay/s6-rc.d/init-chown/run
#!/usr/bin/with-contenv bash

lsiown abc:abc "$HOME"/.cache
