# root/etc/s6-overlay/s6-rc.d/init-chown/run
#!/usr/bin/with-contenv bash

if [ "$MODE" != "local" ]; then
	chown abc:abc -R /app
fi
