# root/etc/s6-overlay/s6-rc.d/file-delivery/run
#!/usr/bin/with-contenv bash

cd /app

if [ -f ".initialized" ]; then
	exit 0
fi

cp -a /defaults/. .
touch .initialized
