# root/etc/s6-overlay/s6-rc.d/file-delivery/run
#!/usr/bin/with-contenv bash

cd /app

if [ -d "app" ]; then
	exit 0
fi

mv /defaults/* .
