# root/etc/s6-overlay/s6-rc.d/file-delivery/run
#!/usr/bin/with-contenv bash

cd /app

if [ -f "requirements.txt" ]; then
	exit 0
fi

cp -r /defaults/* .
